# support/init/akshi-resque.conf
start on starting akshi
stop on stopping akshi
